# derivatives/pytorch/derivatives/stabledaw/ROOT/etc/supervisor/conf.d/stabledaw.conf
[program:stabledaw]
environment=PROC_NAME="%(program_name)s"
command=/opt/supervisor-scripts/stabledaw.sh
autostart=true
autorestart=unexpected
stdout_logfile=/dev/stdout
redirect_stderr=true
